# deploy/sms-printer-ngrok.service
[Unit]
Description=ngrok tunnel for SMS Printer
After=network-online.target sms-printer.service
Wants=network-online.target
Requires=sms-printer.service

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/deploy/ngrok-run.sh
Restart=always
RestartSec=10
StartLimitIntervalSec=600
StartLimitBurst=20
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sms-printer-ngrok

[Install]
WantedBy=multi-user.target
