# deploy/sms-printer-webhook.service
[Unit]
Description=Update Twilio SMS webhook with current ngrok URL
After=sms-printer-ngrok.service
Requires=sms-printer-ngrok.service

[Service]
Type=oneshot
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/venv/bin/python __INSTALL_DIR__/setup_webhook.py
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=600
StartLimitBurst=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sms-printer-webhook

[Install]
WantedBy=multi-user.target
