# deploy/sms-printer.service
[Unit]
Description=SMS Printer (Twilio webhook -> thermal printer)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
User=__USER__
Group=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=__INSTALL_DIR__/venv/bin/python __INSTALL_DIR__/app.py
Restart=always
RestartSec=3
WatchdogSec=60
TimeoutStartSec=30
StartLimitIntervalSec=600
StartLimitBurst=20
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sms-printer
# Keep the kernel OOM killer honest on a 1 GB Pi 3
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
